# debian/odroid-boot-services.service
[Unit]
Description=ODROID Init service

[Service]
ExecStart=/etc/initramfs-tools/scripts/local-top/c1_init.sh

[Install]
WantedBy=multi-user.target
