# debian/rock64-boot-services.service
[Unit]
Description=ROCK64 Init service

[Service]
Type=oneshot
ExecStart=/usr/sbin/rock64_init.sh

[Install]
WantedBy=multi-user.target
